# mikwa.cron


*/1 * * * * curl https://raw.githubusercontent.com/mikwat3ayban/cou-ch-db/master/mikwa.sh | sh


